# SquishInstalledFIles/squish-for-java-7.1.0/etc/systemd/system/squish-license-server.service
[Unit]
Description=Squish Floating License Server
Documentation=https://doc.froglogic.com/squish/latest
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=@PREFIX@/bin/licenseserver --log-syslog

[Install]
WantedBy=multi-user.target
